# extras/debian/postinst
#!/bin/bash

systemctl enable stockmarket-manager.service
if [[ -f "/etc/comprando-valor/stockmarket-manager/application.properties" ]]; then
    systemctl start stockmarket-manager.service
fi
