# extras/debian/postinst
#!/bin/bash

systemctl enable stockmarket-manager.service
systemctl start stockmarket-manager.service
